# supervisord.conf
[program:node]
command = node webserver.js
directory = /home/dotcloud/current
